# router/rootfs/etc/services.d/bind/run
#!/usr/bin/with-contenv sh
#fdmove -c 2 1
#s6-envdir env

echo "Starting bind dns server..."
named -f -c /etc/bind/named.conf
